# imageroot/systemd/user/moodle-smtps.service
[Unit]
Description=SMTP forwarder for Moodle (1587 → mail.ksatdesign.com.au:587)
BindsTo=moodle-app.service
After=moodle-app.service moodle.service moodle-mariadb.service

[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:1587,reuseaddr,fork TCP:mail.ksatdesign.com.au:587
Restart=always

[Install]
WantedBy=multi-user.target
